# systemd/kiosk-backlight.service
[Unit]
Description=Kiosk backlight idle controller
After=graphical-session.target
Wants=graphical-session.target

[Service]
Type=simple
ExecStart=/usr/local/bin/kiosk-backlight.sh
Restart=always
RestartSec=1

[Install]
WantedBy=default.target
