# systemd/kiosk-backlight.service
[Unit]
Description=Kiosk backlight idle controller (evdev + swallow)
After=multi-user.target

[Service]
Type=simple
ExecStart=/usr/local/bin/kiosk-backlight.sh
Restart=always
RestartSec=1

NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=true

[Install]
WantedBy=multi-user.target
